# bin/local.d/nbfc.start
#/bin/bash

sudo /opt/nbfc/nbfcservice.sh start
