# bin/local.d/nbfc.start
#/bin/bash

/opt/nbfc/nbfcservice.sh start
